# workload/mongo-db/4.4.18/conf/supervisor.ini
[program:mongo-db]
command=bash -c "mongod --replSet rs0 --config /opt/mongo-db/conf/server.conf"
process_name=%(program_name)s
pidfile=/opt/mongo-db/system/server.pid
autorestart=true

[program:init-replica]
command=bash -c "/opt/mongo-db/script/init-replica-set.sh"
process_name=%(program_name)s
autorestart=false     ; Disable auto-restart since the script runs only once
startsecs=0           ; Start the script immediately after launching
exitcodes=0           ; Exit code 0 is considered a successful exit
