# workload/mongo-db/4.4.18/conf/supervisor.ini
[program:mongo-db]
command=bash -c "mongod --replSet rs0 --config /opt/mongo-db/conf/server.conf"
process_name=%(program_name)s
pidfile=/opt/mongo-db/system/server.pid
autorestart=true
