# workload/redis/7.5.2/conf/supervisor.ini
[program:redis]
command=bash -c "redis-server /opt/redis/conf/redis.conf"
process_name=%(program_name)s
pidfile=/opt/redis/system/server.pid
autorestart=false
